# supervisor.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/var/run//supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run//supervisor.sock

[program:etherpad]
directory=/opt/etherpad/bin
command=run.sh --root
user=root
autostart=true
autorestart=true
